# server/healthclaw.service
[Unit]
Description=HealthClaw API Server
After=network.target

[Service]
Type=simple
User=lars
WorkingDirectory=/home/lars/.openclaw/workspace-coder/HealthClaw/server
Environment=HEALTHCLAW_API_KEY=hb-lars-2026
Environment=HEALTHCLAW_PORT=8099
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8099
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
